# apps/bridge/2d6mcp-bridge.service
[Unit]
Description=2d6mcp Discord Voice Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bridge
Group=bridge
EnvironmentFile=/opt/2d6mcp/apps/bridge/.env
WorkingDirectory=/opt/2d6mcp
ExecStart=/usr/bin/node /opt/2d6mcp/apps/bridge/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=2d6mcp-bridge

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/2d6mcp
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictNamespaces=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
